# fs/etc/network/if-up.d/iptables
#!/bin/sh
iptables-restore -c < /etc/iptables.rules
